# templates/etc/network/if-pre-up.d/iptablesload
#!/bin/sh
iptables-restore < /etc/sysconfig/iptables.save
exit 0
